# supervisor-app.conf
[program:app-uwsgi]
command=/usr/local/bin/uwsgi --ini /dsdpy-app/department/uwsgi.ini
autostart = true
autorestart = true
startretries = 3
stdout_logfile = /dsdpy-app/department/logs/supervisord.log